# study/foo.service
[Unit]
Description=A Example Systemd Service
Wants=

[Service]
ExecStart=/usr/local/bin/foo
